# deploy/norscode.service
[Unit]
Description=Norscode backend service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/norscode/current
EnvironmentFile=-/etc/norscode/norscode.env
ExecStart=/srv/norscode/current/bin/nc serve app.no --production --host 0.0.0.0 --port 8000 --keep-alive --request-timeout 30 --restart-on-crash --max-restarts 3
Restart=on-failure
RestartSec=2
KillSignal=SIGTERM
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/norscode /var/lib/norscode /var/log/norscode

[Install]
WantedBy=multi-user.target
